# f1-fantasy-league.service
[Unit]
Description=F1 Fantasy Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/f1-fantasy-league
EnvironmentFile=/home/pi/f1-fantasy-league/.env
ExecStart=/home/pi/f1-fantasy-league/venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
